# services/adsb-server.service
[Unit]
Description=ADS-B Server - Data Receiver and Forwarder
After=network.target dump1090-fa.service
Wants=dump1090-fa.service

[Service]
Type=simple
User=jlbmaritime
Group=jlbmaritime
WorkingDirectory=/home/jlbmaritime/adsb-wifi-manager
ExecStart=/usr/bin/python3 /home/jlbmaritime/adsb-wifi-manager/adsb_server/adsb_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
